# corepure64/etc/inittab
# /etc/inittab: init configuration for busybox init - Unikernel optimized
# Boot-time system configuration/initialization script.
#
::sysinit:/etc/init.d/rcS

# Single TTY for unikernel application
tty1::respawn:/sbin/getty 38400 tty1

# Shutdown handling
::restart:/etc/init.d/rc.shutdown
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rc.shutdown